# supervisor/kafka.conf
[program:kafka]
command=/usr/bin/start-kafka.sh
autostart=true
autorestart=true

[program:create-topics]
command=/usr/bin/create-topics.sh
startsecs = 10
autostart=true
autorestart=false
